# motd/update-motd.d/99-tb_init_bootmenu-notice
#!/bin/sh
echo "\r"
echo "%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%"
echo "% Execute \e[30;38;5;215mtb_init_bootmenu\e[0;0m to configure the \e[30;38;5;215moverlay \e[30;38;5;130mmode\e[0;0m/\e[30;38;5;130moptions\e[0;0m %"
echo "%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%%"
echo "\r"
